# files/etc/crontabs/root
# Check for new firmware every 6 hours
0 */6 * * * /usr/sbin/autoupdater
# Publish map data
*/5 * * * * /usr/sbin/print_map.sh -p
# Publish service data
*/5 * * * * /usr/sbin/print_service.sh -p
# Reboot every 5 days at 04:05
5 4 */5 * * /sbin/reboot
