# files/etc/crontabs/root
#check gateway and WAN connectivity every 5 minutes
*/5 * * * * /etc/init.d/freifunk_watchdog start
#check for new firmware once a day
0 * * * * /usr/sbin/autoupdater
#reboot every 5 days at 04:05
5 4 */5 * * /sbin/reboot
