# files/etc/crontabs/root
*/5 * * * * /sbin/wan_watchdog > /dev/null