# files/etc/crontabs/root
*/5 * * * * /etc/init.d/freifunk_watchdog start
#reboot every 5 days at 04:05
5 4 */5 * * /sbin/reboot
