# files/etc/hotplug.d/net/11-batman
#this script receives the mac of the best gateway in reach
#from batman-adv and sets the routing table accordingly

if [ "$ACTION" = "change" ]; then

  exec >/tmp/11-batman.log 2>&1
  date

  if [ "$BATACTION" = "del" ]; then
    /sbin/set_gateway "mesh" ""
    exit 0
  fi

  if [ -n "$BATDATA" ]; then
    echo "(I) 11-batman: received gateway MAC: '$BATDATA'"
    GATE_IP=`mac2ip $BATDATA`
    /sbin/set_gateway "mesh" "$GATE_IP"
    exit 0
  fi
fi
