# files/etc/hotplug.d/net/11-batman
#receive the mac of the best gateway in reach from batman-adv

[ "$ACTION" != "change" ] && exit

/etc/init.d/freifunk_watchdog start
