# files/etc/hotplug.d/net/11-batman
exit 0

if [ "$ACTION" = "change" ] ; then

 if [ "$BATACTION" = "del" ] ; then
    echo "del default route" >> /dev/kmsg
    route del default
    exit 0
 fi

 if [ -n "$BATDATA" ] ; then
    echo "$BATDATA" > /tmp/11-batman.log
    GATE_IP=`mac2ip $BATDATA`
    if [ $? -eq 1 -o -z "$GATE_IP" ]; then
        exit 1
    fi
    echo "update bat0 route to $BATDATA ($GATE_IP)" >> /dev/kmsg
    uci set dhcp.@dnsmasq[0].dhcp_option="option:router, $GATE_IP"
    route del default
    route add default gw "$GATE_IP"
    /etc/init.d/dnsmasq reload
    echo "search lan" > /etc/resolv.conf
    echo "nameserver $GATE_IP" >> /etc/resolv.conf
  fi
fi
