# package/freifunk-base/files/etc/uci-defaults/90_freifunk_install_firewall.sh
#!/bin/sh

cp /usr/share/freifunk-base/firewall.user /etc/

exit 0
