# package/freifunk-base/files/etc/uci-defaults/99_installcron
#!/bin/bash

cat /lib/cron/* /etc/crontabs/root

#execute on every restart before initscripts
exit 1
