# files/etc/init.d/freifunk_watchdog
#!/bin/sh /etc/rc.common

#This script is run at startup and then every 5min using cron.

START=99

start()
{

#Check for Internet access on the WAN port and switch
#batman-adv between client und server mode accordingly.
#This script is called by cron every 5min and by /etc/hotplug.d/net/11-batman.

#prevent script from being called simultaneously
lock_dir=/var/lock/freifunk_watchdog.lock
mkdir $lock_dir || exit

exec >/tmp/freifunk_watchdog.log 2>&1


echo "(I) Start: $(date)"

[ "$(uci get -q firewall.share_internet)" = "forwarding" ] && \
	share_internet="yes" || share_internet="no"
echo "(I) share_internet: '$share_internet'"

addr() {
	ip -$1 address show dev $2 2> /dev/null | tr '/' ' '| awk '/inet/{print($2); exit(0);}'
}

default_gw() {
	ip -$1 route list 0/0 dev $2 2> /dev/null | awk '{print($3); exit(0)}'
}

own_public_ip=$(addr 4 br-public)
own_private_ip=$(addr 4 br-private)

echo "(I) Own public IPv4 address: '$own_public_ip'"
echo "(I) Own private IPv4 address: '$own_private_ip'"

fastd=$(which fastd)


#shuffle a list and pick n entries
pick()
{
	local n="$1"; shift;
	echo $@ | tr ' ' '\n' | awk 'BEGIN { srand() } { print rand() "\t" $0 }' | sort -n | cut -f2- | head -$n
}

#ping root servers to check a gateway
#on a given interface for internet access
check_internet() {
	local dev="$1" gw="$2"
	local addrs="192.33.4.12 128.8.10.90 193.0.14.129 198.41.0.4 192.228.79.201 192.5.5.241 192.36.148.17 192.58.128.30"
	local addr cmd ok

	[ -z "$gw" ] && { return 1; }
	[ -z "$dev" ] && { echo "(E) No WAN interface given for check_internet()."; return 1; }

	echo "(I) Try to ping root servers on '$dev' using gateway: '$gw'"
	IFS="
"
	for addr in `pick 3 $addrs`; do
		echo -n "(I) Ping $addr "
		ip route add "$addr/32" via "$gw" dev "$dev"
		ping "$addr" -c 1 -I "$dev" -q -W 1 &> /dev/null; ok=$?;
		ip route del "$addr/32" via "$gw" dev "$dev"
		if [ $ok -eq 0 ]; then
			echo "- success."
			return 0
		else
			echo "- failed."
		fi
	done
	echo "(W) Cannot reach any tested IP."
	return 1
}

#check/establish internet on wan
check_wan()
{
	local wan="$1"
	local gw=""

	[ -z "$wan" ] && {
		echo "(I) No WAN interface found."
		return 1
	}

	echo "(I) Check if Internet is accessible via '$wan'."

	gw=$(default_gw 4 "$wan")
	check_internet "$wan" "$gw" && return 0

	[ "$(uci get -q network.wan.proto)" != "dhcp" ] && {
		return 0
	}

	echo "(I) Try to aquire new DHCP lease."
	local pidfile="/var/run/udhcpc-$wan.pid"
	if [ -f "$pidfile" ]; then
		local pid=$(cat "$pidfile")
		kill -s SIGUSR2 "$pid" #release lease
		kill -s SIGUSR1 "$pid" #get new lease

		sleep 3
		gw=$(default_gw 4 "$wan")
		check_internet "$wan" "$gw" && return 0
	else
		echo "(W) No pidfile found."
	fi

	echo "(W) No Internet detected on '$wan'"
	ip addr flush dev "$wan" 2> /dev/null
	return 1
}

is_running() {
	ps | grep -v grep | grep "$1" > /dev/null
}

#Get the best public gateway.
find_public_gateway() {
	echo "(I) Try to get the address of a public gateway."
	IFS="
"
	#this disables the gateway check for uptime < 5min
	skip=$(uptime | awk '{printf("%d",$4=="min," && $3 < 5)}')

	#TODO: ask alfred to sort the resutls for us
	for addr in $(alfred -r 88 | tr '{}[]\n"\\,' ' ' | awk 'BEGIN{c=0}{ if($4 == "gateway" && $7 ~ /^10.[0-9.]{5,11}$/) {print($7); if(c>9) exit(0)}; c++;}' 2> /dev/null); do
		if check_internet "br-public" "$addr" || test $skip = "1"; then
			public_gw_mac="$mac"
			public_gw_ip="$addr"
			echo "(I) Best public gateway: '$public_gw_ip'"
			return;
		fi
	done
}

#Configure the DHCP server and routing table.
set_routes()
{
	local has_internet="$1"
	local public_str private_str

	echo "(I) Start setting routes:"

	[ -z "$own_private_ip" ] && { echo "(E) Private IP address not set."; return; }
	[ -z "$own_public_ip" ] && { echo "(E) Public IP address not set."; return; }

	#DHCP settings for all states: ('3' for gateway, '6' for DNS-server)
	if [ -n "$public_gw_ip" ]; then
		if [ "$has_internet" = "yes" ]; then
			if [ "$share_internet" = "yes" ]; then
				public_str="3,$own_public_ip"
				private_str="3,$own_private_ip"
			else
				public_str="3,$public_gw_ip 6,$public_gw_ip"
				private_str="3,$own_private_ip"
			fi
		else
			public_str="3,$public_gw_ip 6,$public_gw_ip"
			private_str="3,$own_private_ip 6,$public_gw_ip"
		fi
	else
		if [ "$has_internet" = "yes" ]; then
			if [ "$share_internet" = "yes" ]; then
				public_str="3,$own_public_ip 6,$own_public_ip"
				private_str="3,$own_private_ip 6,$own_public_ip"
			else
				public_str="3"
				private_str="3,$own_private_ip"
			fi
		else
			public_str="3"
			private_str="3,$own_private_ip"
		fi
	fi

	if [ "$has_internet" != "yes" ]; then
		echo "(I) Clear default routes and nameserver list."
		while ip route del default 2> /dev/null; do true; done

		[ -n "$public_gw_ip" ] && {
			echo "(I) Set default route and nameserver to '$public_gw_ip'."
			ip route add default via "$public_gw_ip"
		}
	fi

	#entries for the Intercity-VPN network
	for net in 10.0.0.0/8 172.16.0.0/12; do
		if [ -n "$public_gw_ip" ]; then
			local EXIST=$(ip route show "$net" | grep -c -m 1 "$public_gw_ip" 2> /dev/null)
			if [ "$EXIST" -eq 0 ]; then
				ip route del "$net" 2> /dev/null
				ip route add "$net" via "$public_gw_ip"
			fi
		else
			ip route del "$net" 2> /dev/null
		fi
	done

	local pre_public_str="$(uci get dhcp.public.dhcp_option)"
	local pre_private_str="$(uci get dhcp.private.dhcp_option)"

	if [ "$public_str" != "$pre_public_str" -o "$private_str" != "$pre_private_str" ]
	then
		uci set dhcp.public.dhcp_option="$public_str"
		uci set dhcp.private.dhcp_option="$private_str"

		echo "(I) DHCP settings were changed => reload dnsmasq."
		/etc/init.d/dnsmasq reload

	else
		echo "(I) No change of DHCP settings."
	fi

	echo "(I) DHCP private options: '$private_str'"
	echo "(I) DHCP public option: '$public_str'"
}

if check_wan "br-wan"; then
	echo "(I) We have local Internet."
	has_internet="yes"
else
	echo "(I) We do not have local Internet."
	has_internet="no"
fi

#Get the best public gateway in the public network.
find_public_gateway

#store data for web-gui
echo -n "$has_internet" > /tmp/has_internet


if [ "$has_internet" = "yes" ]
then
	set_routes "yes"

	#start/restart fastd
	is_running "$fastd" || /etc/init.d/fastd start

	#announce gateway
	[ "$share_internet" = "yes" ] && {
		echo -n "{\"type\" : \"gateway\", \"addr\" : \"$own_public_ip\"}" | alfred -s 88
	}
else
	set_routes "no"

	#stop fastd
	is_running "$fastd" && /etc/init.d/fastd stop
fi


echo "(I) Done."
rmdir $lock_dir
}
