# files/etc/init.d/freifunk_first
#!/bin/sh /etc/rc.common

START=31


start () 
{
exec >/tmp/freifunk.log 2>&1

echo "(I) Freifunk: start freifunk_first"

mac=`uci get freifunk.settings.mac`
if [ -n "$mac" ]; then
    
    echo "(I) Freifunk: Use mac: '$mac'"

    #Create dummy interfaces to be able to set the mac of
    #br-mesh and the mac for the OGMs of batman-adv freely.

    echo "(I) Freifunk: Create dummy_ext interface."
    tunctl -t dummy_ext > /dev/null
    ifconfig dummy_ext hw ether $mac up

    echo "(I) Freifunk: Create dummy_mesh interface."
    tunctl -t dummy_mesh > /dev/null
    ifconfig dummy_mesh hw ether $mac up
fi

echo "(I) Freifunk: done freifunk_first"
}
