# files/etc/hotplug.d/iface/22-wan_led
#!/bin/sh
case "${ACTION:-ifup}" in
        ifup)
		gpioctl dirout 1
                gpioctl set 1
        ;;
        ifdown)
                gpioctl clear 1
        ;;
esac
