# files/etc/hotplug.d/iface/30-fastd
#!/bin/sh


[ "$ACTION" = "ifup" -a "$INTERFACE" = "wan" ] || exit 0

uci set fastd.default.bind="any interface \"$DEVICE\""

/etc/init.d/fastd start
