# files/etc/hotplug.d/iface/40-freifunk
#!/bin/sh

[ "$INTERFACE" = "public6" ] || exit
  
. /lib/functions/network.sh

dnsserver=""
network_get_dnsserver dnsserver public6
echo "nameserver $dnsserver" > /tmp/resolv.conf.public.auto

#/etc/init.d/dnsmasq2 reload
